# tasmoadmin/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TasmoAdmin
# Runs the PHP-FPM daemon
# ==============================================================================
bashio::log.info "Starting PHP-FPM server..."
exec php-fpm81 -R --nodaemonize
