# tasmoadmin/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: TasmoAdmin
# Runs the Nginx daemon
# ==============================================================================

# Wait for PHP-FPM to become available
bashio::net.wait_for 9001

bashio::log.info "Starting NGINX server..."
exec nginx -g "daemon off;"
