# tasmoadmin/rootfs/etc/services.d/php-fpm/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: TasmoAdmin
# Runs the PHP-FPM daemon
# ==============================================================================
bashio::log.info "Starting PHP-FPM server..."
exec php-fpm8 -R --nodaemonize
