# tasmoadmin/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: TasmoAdmin
# Runs the PHP-FPM daemon
# ==============================================================================
exec php-fpm7 --nodaemonize
